# deployment/daddeals.service.example
# DadDeals web app systemd service example.
#
# Raspberry Pi path for Nolan's DadDeals install:
#   /home/nolan/DadDeals
#
# If your Pi username or project folder changes later, update the User,
# Group, WorkingDirectory, EnvironmentFile, and ExecStart paths below.

[Unit]
Description=DadDeals Flask web app
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nolan
Group=nolan
WorkingDirectory=/home/nolan/DadDeals
EnvironmentFile=-/home/nolan/DadDeals/.env
ExecStart=/home/nolan/DadDeals/.venv/bin/gunicorn --workers 1 --threads 2 --bind 0.0.0.0:5000 app:app
Restart=on-failure
RestartSec=5

# Keep logs in the normal systemd journal. The service does not print secrets,
# and .env values are loaded from the file instead of being written here.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
